# systemd/oeil-api.service
[Unit]
Description=Oeil — FastAPI backend (cameras, recording, API)
After=network-online.target oeil-go2rtc.service
Wants=network-online.target
Requires=oeil-go2rtc.service

[Service]
Type=simple
User=oeil
Group=oeil
WorkingDirectory=/opt/oeil
EnvironmentFile=/etc/oeil/oeil.env
ExecStart=/opt/oeil/venv/bin/uvicorn main:app \
    --host ${OW_HOST} \
    --port ${OW_PORT} \
    --workers 2 \
    --log-level info \
    --access-log
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=oeil-api

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/oeil /var/log/oeil /etc/oeil
ProtectHome=true

[Install]
WantedBy=multi-user.target
