# config/oeil-go2rtc.service
# ── oeil-go2rtc.service ──────────────────────────────────────────────────
# RTSP/ONVIF proxy + HLS/WebRTC bridge

[Unit]
Description=Oeil — go2rtc stream proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=oeil
Group=oeil
ExecStart=/opt/oeil/scripts/start_go2rtc.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=oeil-go2rtc

[Install]
WantedBy=multi-user.target
