# config/oeil-go2rtc.service
[Unit]
Description=Oeil — go2rtc stream proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=oeil
Group=oeil
LoadCredentialEncrypted=CAM_PASS:/etc/credstore.encrypted/CAM_PASS
ExecStart=/bin/bash -c 'export CAM_USER=admin; export CAM_PASS=$(cat $CREDENTIALS_DIRECTORY/CAM_PASS); exec /usr/local/bin/go2rtc -config /etc/oeil/go2rtc.yaml'
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=oeil-go2rtc

[Install]
WantedBy=multi-user.target
